# sh/low-battery.service
[Unit]
Description=Check for low battery

[Service]
Type=simple
ExecStart=/usr/local/bin/lowbatt
PIDFile=/tmp/lowbatt.pid
